# deb/prerm
#!/bin/bash

systemctl stop enunu
systemctl disable enunu
rm /lib/systemd/system/enunu.service